# packaging/com.imessage.ai.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.imessage.ai</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>BOT_DISABLE_MONITOR</key>
        <string>__BOT_DISABLE_MONITOR__</string>
        <key>IMSG_AI_TOKEN</key>
        <string>__IMSG_AI_TOKEN__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/Library/Logs/imessage-ai.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/imessage-ai.log</string>
</dict>
</plist>
